# config/supervisor/psql.conf
[program:psql]
command=/usr/lib/postgresql/9.5/bin/postgres -D /data/psql -c config_file=/etc/postgresql/9.5/main/postgresql.conf
user=postgres
autorestart=true

